# home/.config/systemd/user/gpg-agent.service.d/override.conf
# disable systemd --user controlled gpg-agent. Gpg agent is managed by gpg
# itself with these dotfiles.
[Service]
ExecStart=
ExecStart=/bin/true
ExecReload=
ExecReload=/bin/true
