# home/.config/systemd/user/gpg-agent.socket.d/override.conf
# disable systemd --user controlled gpg-agent. Gpg agent is managed by gpg
# itself with these dotfiles.
[Socket]
ListenStream=%t/gnupg/S.gpg-agent.disabled
