# etc/update-motd.d/00-header
#!/bin/sh
printf "\033c"
echo
figlet $(hostname -s)
echo "     .$(hostname -d)"
