# etc/update-motd.d/00-header
#!/bin/sh

# clear terminal
printf "\033c"
echo

# moar ascii art!
figlet $(hostname -s)
echo "     .$(hostname -d)"
